# trunk/other/movie/crawler/crawler.cron
30 10 * * * /home/movie/cron/crawler.sh>>/home/movie/cron/crawler.log 2>>/home/movie/cron/crawler.log
