# other/movie/notify/movie_notify.cron
*/30 * * * * /home/cbib/news/news.py>>/home/cbib/news/news.log 2>>/home/cbib/news/news.log